# site/minecraft/templates/minecraft.service.epp
[Unit]
Description=Minecraft Server

Wants=network.taget
After=network.taget

[Service]
WorkingDirectory=<%= $instal_dir %>
ExecStart=/usr/bin/java -Xmx512M -Xms32M -jar server.jar nogui

[Install]
WantedBy=multi-user.target
